# apps/migrapilot-console/ops/systemd/migrapilot-inventory-refresh.timer
[Unit]
Description=Run MigraPilot inventory refresh every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=migrapilot-inventory-refresh.service

[Install]
WantedBy=timers.target
